# nginx-certbot/cron/crontab
# m h  dom mon dow   command
0 3  * * *  ~/cron_job.sh
